# snort/templates/snort_upstart.conf
description "Snort network intrusion detection"
stop on runlevel [!2345]
start on runlevel [2345]
script
    exec /usr/local/bin/snort -c {{ snort_config_path }}
end script
